# docker/etc/s6/jira/run
#!/bin/bash

# 配置
source ./setup

# 启动Jira服务器
exec gosu "${USERNAME}" /opt/atlassian/jira/bin/start-jira.sh -fg
